# pi2/services/playback.service
[Unit]
Description=Timelapse Playback Service
After=sync.service
Wants=sync.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/playback.py
WorkingDirectory=/data
Restart=always
RestartSec=10
User=root
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1
# DRM/KMS access for mpv on Waveshare round display
Environment=XDG_RUNTIME_DIR=/run/user/0
# Grant access to DRM/graphics devices
SupplementaryGroups=video render

[Install]
WantedBy=multi-user.target
