# pi2/services/status_api.service
[Unit]
Description=Timelapse Status API (Flask)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/status_api.py
WorkingDirectory=/data
Restart=always
RestartSec=10
User=root
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1
# On Bookworm, ensure dual-stack sockets can accept both IPv4 and IPv6
Environment=FLASK_RUN_HOST=::

[Install]
WantedBy=multi-user.target
