# pi2/services/sync.service
[Unit]
Description=Timelapse Image Sync Service
After=network-online.target remote-fs.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/sync.py
WorkingDirectory=/data
Restart=always
RestartSec=10
User=root
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1
WatchdogSec=180
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target
