# pi2/backup.service
[Unit]
Description=Pi2 Daily Timelapse Backup
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/opt/timelapse/backup.sh
StandardOutput=journal
StandardError=journal
