# pi1/services/portal.service
[Unit]
Description=Timelapse Admin Portal (Flask)
After=network.target local-fs.target
Wants=network.target
RequiresMountsFor=/data

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/portal.py
WorkingDirectory=/data
Restart=always
RestartSec=10
User=root
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
