# pi1/services/wlan0-static-ip.service
[Unit]
Description=Assign static IP 192.168.50.1/24 to wlan0 (Timelapse AP)
# Run after hostapd has started so the wlan0 interface exists in AP mode.
# Works regardless of whether dhcpcd or NetworkManager manages networking,
# because both may leave wlan0 without an address when NM is told to ignore it.
After=hostapd.service
Wants=hostapd.service
# dnsmasq must not start until the IP is on the interface
Before=dnsmasq.service

[Service]
Type=oneshot
RemainAfterExit=yes
# Bring the interface up (idempotent), flush any existing IPv4 addresses,
# then assign the static AP gateway address.
ExecStart=/bin/bash -c '\
  ip link set wlan0 up; \
  ip -4 addr flush dev wlan0 2>/dev/null || true; \
  ip addr add 192.168.50.1/24 dev wlan0 2>/dev/null || true'

[Install]
WantedBy=multi-user.target
